# examples/airtouch5-webui.service
# airtouch5-webui.service
#
# Example systemd unit for the AirTouch 5 web UI. Copy to
# /etc/systemd/system/ (or ~/.config/systemd/user/ for a user service) and
# adjust the ExecStart path / env vars to taste.
#
# Install:
#   1. Build and install the binary:
#        cargo build --release
#        sudo install -m 0755 target/release/airtouch5-webui \
#          /usr/local/bin/airtouch5-webui
#   2. Create the service user (static, so its state dir is stable across
#      reboots -- needed to persist the automation config):
#        sudo useradd --system --no-create-home --shell /usr/sbin/nologin \
#          airtouch5
#   3. Copy this unit file:
#        sudo install -m 0644 examples/airtouch5-webui.service \
#          /etc/systemd/system/airtouch5-webui.service
#   4. Reload systemd and enable + start it:
#        sudo systemctl daemon-reload
#        sudo systemctl enable --now airtouch5-webui
#   5. Check it is running and view logs:
#        systemctl status airtouch5-webui
#        journalctl -u airtouch5-webui -f
#
# Configuration: all options are CLI flags (no env vars). Pass them on the
# ExecStart line; the defaults match the binary's own defaults (bind
# 0.0.0.0:3000, 3s discovery, 60s automation tick). Logging is the one
# env-driven option and stays on RUST_LOG below.
#   --bind <addr:port>              -- bind addr (default 0.0.0.0:3000)
#   --discovery-timeout-ms <ms>     -- UDP discovery wait (default 3000)
#   --automation-tick-secs <s>     -- automation eval tick (default 60; 0 = off)
#   --automation-config <path>     -- automation config file path
#   RUST_LOG                        -- tracing filter (default ...=info,tower_http=info)

[Unit]
Description=AirTouch 5 web UI
Documentation=https://github.com/mason-larobina/airtouch5-webui
# The app discovers the console via UDP, so it needs the network up.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=airtouch5
Group=airtouch5

# /var/lib/airtouch5 is created by systemd and owned by the service user.
# The automation config is persisted here so toggles survive restarts.
StateDirectory=airtouch5

# --- Runtime options (tweak as needed) -------------------------------------
# All options are CLI flags. Adjust the flags below to taste.
ExecStart=/usr/local/bin/airtouch5-webui --automation-config /var/lib/airtouch5/automation.json

# Logging is the one env-driven option: RUST_LOG sets the tracing filter
# (default airtouch5_webui=info,tower_http=info).
Environment=RUST_LOG=airtouch5_webui=info,tower_http=info

# --- Restart policy --------------------------------------------------------
Restart=on-failure
RestartSec=5

# --- Hardening -------------------------------------------------------------
# No capabilities needed: the app listens on port 3000 (>1024) and does UDP
# discovery + a TCP console connection, all as an unprivileged user.
NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=

# Read-only filesystem except the state dir (handled by StateDirectory above)
# and /tmp / /var/tmp (PrivateTmp gives a private, writable /tmp).
ProtectSystem=strict
PrivateTmp=yes
ReadWritePaths=/var/lib/airtouch5

# Network: keep real access (the app does UDP discovery), but hide /home and
# other users' processes.
ProtectHome=yes
ProtectProc=invisible
PrivateDevices=yes

# Kernel protections.
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
ProtectKernelLogs=yes

# Restrict to the syscall set a typical service needs; deny the rest.
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RemoveIPC=yes

# Only the address families the app uses:
#   AF_INET / AF_INET6 -- UDP broadcast discovery, the TCP console connection,
#                          and the HTTP listener.
#   AF_UNIX            -- tokio / tracing internals, /tmp bits.
#   AF_NETLINK         -- required by libc::getifaddrs(), which the airtouch5
#                          discovery crate calls to enumerate network
#                          interfaces and derive their broadcast addresses.
#                          Without this, RestrictAddressFamilies blocks the
#                          netlink socket with EAFNOSUPPORT and discovery
#                          cannot find the console on the local network.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK

[Install]
WantedBy=multi-user.target
